# pkg/deb/control/prerm
#!/bin/sh -e
if [ -x "/etc/init.d/disthene-reader" ]; then
	invoke-rc.d disthene-reader stop || exit $?
fi
